# zohomap.service
[Unit]
Description=Gunicorn instance to serve ZohoMap
After=network.target zohomap-secrets.service
Requires=zohomap-secrets.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/zohomap
EnvironmentFile=/etc/zohomap/app.env
ExecStart=/var/www/zohomap/venv/bin/gunicorn --workers 3 --worker-class gthread --threads 4 --timeout 120 --bind unix:zohomap.sock -m 007 app:app
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
