# zohomap-sync.timer
[Unit]
Description=ZohoMap Auto-Sync Scheduler Timer

[Timer]
# Run every 5 minutes to check if auto-sync is due
OnCalendar=*:0/5
Persistent=true

[Install]
WantedBy=timers.target
