# zohomap-secrets.service
[Unit]
Description=ZohoMap AWS Secrets Loader
Documentation=Fetches secrets from AWS Secrets Manager and writes /etc/zohomap/app.env
After=network-online.target
Wants=network-online.target
# Run before the main app and nightly sync
Before=zohomap.service zohomap-sync.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
ExecStart=/var/www/zohomap/load_aws_secrets.sh
StandardOutput=journal
StandardError=journal
# Retry up to 3 times if AWS is temporarily unavailable at boot
Restart=on-failure
RestartSec=10s
StartLimitBurst=3
StartLimitIntervalSec=60s

[Install]
WantedBy=multi-user.target
